# com.bard.brain-mcp-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bard.brain-mcp-server</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/bard/Code/claude-brain/index.js</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/bard/Code/claude-brain</string>
    
    <key>StandardOutPath</key>
    <string>/Users/bard/Library/Logs/Claude/brain-mcp-server.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/bard/Library/Logs/Claude/brain-mcp-server.error.log</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>
</dict>
</plist>